# fastapi.service
[Unit]
Description=Boot up fastapi

[Service]
ExecStart=/bin/bash /home/FastAPI_MongoDB/start.sh

[Install]
WantedBy=multi-user.target